# root/etc/services.d/alerthub/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc node /app/src/index.mjs